# files/tt-rss.service.in
[Unit]
Description=tt-rss uwsgi service
Requires=tt-rss.socket
After=tt-rss.socket
Before=tt-rss-update.service

[Service]
Type=notify
DynamicUser=yes
User=tt-rss
StateDirectory=tt-rss
EnvironmentFile=/etc/tt-rss/config.env
ExecStartPre=@coreutils@/bin/mkdir -p \
        /var/lib/tt-rss/lock \
        /var/lib/tt-rss/cache/images \
        /var/lib/tt-rss/cache/upload \
        /var/lib/tt-rss/cache/export \
        /var/lib/tt-rss/feed-icons
ExecStart=@uwsgi@/bin/uwsgi --ini /tt-rss.ini
Restart=always
MemoryDenyWriteExecute=no
BindReadOnlyPaths=/etc/ssl/certs/

[Install]
WantedBy=multi-user.target
Also=tt-rss.socket
