# files/tt-rss-update.service.in
[Unit]
Description=tt-rss update service
After=tt-rss.service

[Service]
Type=simple
DynamicUser=yes
User=tt-rss
ConfigurationDirectory=tt-rss
StateDirectory=tt-rss
ExecStart=@php@/bin/php @ttRss@/update_daemon2.php
MemoryDenyWriteExecute=no
BindReadOnlyPaths=/etc/ssl/certs/

[Install]
WantedBy=multi-user.target
