# files/tt-rss-update.service.in
[Unit]
Description=tt-rss update service
After=tt-rss.service

[Service]
Type=simple
User=tt-rss
EnvironmentFile=/etc/tt-rss/config.env
StateDirectory=tt-rss
ExecStart=@php@/bin/php @ttRss@/update_daemon2.php
Restart=always

[Install]
WantedBy=multi-user.target
