# files/tt-rss-update.service.in
[Unit]
Description=tt-rss update service
After=tt-rss.service

[Service]
Type=oneshot
DynamicUser=yes
User=tt-rss
StateDirectory=tt-rss
EnvironmentFile=/etc/tt-rss/config.env
ExecStart=@php@/bin/php @ttRss@/update.php --feeds
MemoryDenyWriteExecute=no
BindReadOnlyPaths=/etc/ssl/certs/

[Install]
WantedBy=multi-user.target
