# files/tt-rss.service.in
[Unit]
Description=tt-rss uwsgi service
Requires=tt-rss.socket
After=tt-rss.socket
Before=tt-rss-update.service

[Service]
Type=notify
DynamicUser=yes
User=tt-rss

TemporaryFileSystem=/var/lib
StateDirectory=tt-rss
Environment=SSL_CERT_FILE=/etc/ssl/certs/ca-bundle.crt
Environment=TTRSS_ICONS_DIR=/var/lib/tt-rss/feed-icons
Environment=TTRSS_CACHE_DIR=/var/lib/tt-rss/cache
Environment=TTRSS_LOCK_DIRECTORY=/var/lib/tt-rss/lock
EnvironmentFile=/etc/default/tt-rss.conf

ExecStartPre=@coreutils@/bin/mkdir -p \
        /var/lib/tt-rss/lock \
        /var/lib/tt-rss/cache/images \
        /var/lib/tt-rss/cache/upload \
        /var/lib/tt-rss/cache/export \
        /var/lib/tt-rss/feed-icons

ExecStartPre=-@php@/bin/php @ttRss@/update.php --update-schema=force-yes
ExecStart=@uwsgi@/bin/uwsgi --ini @uwsgiConfig@
Restart=always

[Install]
WantedBy=multi-user.target
Also=tt-rss.socket
