# files/tt-rss-update.service.in
[Unit]
Description=tt-rss update service
After=tt-rss.service

[Service]
Type=simple
DynamicUser=yes
User=tt-rss

TemporaryFileSystem=/var/lib
StateDirectory=tt-rss
Environment=TTRSS_ICONS_DIR=/var/lib/tt-rss/feed-icons
Environment=TTRSS_CACHE_DIR=/var/lib/tt-rss/cache
Environment=TTRSS_LOCK_DIRECTORY=/var/lib/tt-rss/lock
EnvironmentFile=/etc/default/tt-rss.conf

ExecStart=@php@/bin/php @ttRss@/update_daemon2.php
Restart=on-failure
SyslogIdentifier=%N

[Install]
WantedBy=multi-user.target
